# dnfdaemon-server/dbus/dnfdaemon.service
[Unit]
Description=DNF daemon

[Service]
Type=dbus
BusName=org.rpm.dnf.v0
User=root
ExecStart=/usr/bin/dnfdaemon-server
